# Ansible/my_app.service.j2
[Unit]
Description=Your Application Service
After=network.target

    [Service]
ExecStart=ExecStart=/home/mkurmanalieva/myapp/main
WorkingDirectory=/home/mkurmanalieva/myapp

User=meerimkurmanalieva
Restart=always

    [Install]
WantedBy=multi-user.target

